# mdnsd.service.in
[Unit]
Description=mDNS Daemon
After=network.target

[Service]
ExecStart=@SBINDIR@/mdnsd
User=nobody

[Install]
WantedBy=multi-user.target
